# scripts/librix/ldc-daemon
#!/sbin/runscript

depend() {
  need net
  before xdm
}

start() {
  ebegin "Starting LDC"
  start-stop-daemon --start --exec /usr/bin/ldc
  einfo "Finishing LDC"
  eend $?
}

